# deployment/systemd/skyguard.service
[Unit]
Description=SkyGuard Raptor Detection System
After=network-online.target
Wants=network-online.target
Before=skyguard-web.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/skyguard
Environment=PATH=/home/pi/skyguard/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/pi/skyguard/venv/bin/python -m skyguard.main --config config/skyguard.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skyguard

# Health check - verify camera snapshot is being created
ExecStartPost=/bin/sleep 5
ExecStartPost=/usr/bin/test -f /home/pi/skyguard/data/camera_snapshot.jpg

[Install]
WantedBy=multi-user.target

